# tools/hotplug/Linux/systemd/proc-xen.mount.in
[Unit]
Description=Mount /proc/xen files
ConditionVirtualization=xen
RefuseManualStop=true

[Mount]
What=xenfs
Where=/proc/xen
Type=xenfs
