# tools/hotplug/Linux/systemd/var-lib-xenstored.mount.in
[Unit]
Description=mount xenstore file system
Requires=proc-xen.mount
After=proc-xen.mount
ConditionPathExists=/proc/xen/capabilities
RefuseManualStop=true

[Mount]
What=xenstore
Where=@XEN_LIB_STORED@
Type=tmpfs
Options=mode=755
