# tools/hotplug/Linux/systemd/xen-watchdog.service.in
[Unit]
Description=Xen-watchdog - run xen watchdog daemon
Requires=proc-xen.mount
After=proc-xen.mount xendomains.service
ConditionPathExists=/proc/xen/capabilities

[Service]
Type=forking
ExecStart=@sbindir@/xenwatchdogd 30 15
KillSignal=USR1

[Install]
WantedBy=multi-user.target
