# tools/hotplug/Linux/systemd/xendomains.service.in
[Unit]
Description=Xendomains - start and stop guests on boot and shutdown
Requires=xenstored.socket
After=xenstored.service xenconsoled.service
ConditionPathExists=/proc/xen/capabilities

[Service]
Type=oneshot
RemainAfterExit=true
ExecStartPre=/bin/grep -q control_d /proc/xen/capabilities
ExecStart=-@LIBEXEC_BIN@/xendomains start
ExecStop=@LIBEXEC_BIN@/xendomains stop
ExecReload=@LIBEXEC_BIN@/xendomains restart

[Install]
WantedBy=multi-user.target
