# tools/hotplug/Linux/systemd/xenstored_ro.socket.in
[Unit]
Description=xenstore ro socket
ConditionPathExists=/proc/xen/capabilities

[Socket]
ListenStream=/var/run/xenstored/socket_ro
SocketMode=0660
Service=xenstored.service

[Install]
WantedBy=sockets.target
