# tools/hotplug/Linux/systemd/xenstored_ro.socket.in
[Unit]
Description=xenstore ro socket
Requires=proc-xen.mount var-lib-xenstored.mount
After=proc-xen.mount var-lib-xenstored.mount
ConditionPathExists=/proc/xen/capabilities

[Socket]
ListenStream=/var/run/xenstored/socket_ro
SocketMode=0660
Service=xenstored.service

[Install]
WantedBy=sockets.target
