# tools/xm-test/ramdisk/skel/etc/inittab
# This is run first except when booting in single-user mode.
#
::sysinit:/etc/init.d/rcS

::askfirst:-/bin/sh
